# inittab
id:2:initdefault:
si::sysinit:/etc/init.d/rcS
1:2345:respawn:/sbin/getty 38400 tty1
2:23:respawn:/sbin/getty 38400 tty2
3:23:respawn:/sbin/getty 38400 tty3
4:23:respawn:/sbin/getty 38400 tty4
5:23:respawn:/sbin/getty 38400 tty5
ca::ctrlaltdel:/etc/init.d/reboot

